# myapp_local.conf
[program:myserver]
command=gunicorn --log-level info views:app -w 4 -b 0.0.0.0:5000

[supervisord]
logfile=/media/kristy/UltraFit/Dropbox/KristyInsight/app/supervisord.log
loglevel=debug
user=root